# raxdns.cron
*/5 * * * * root php -q /opt/raxdns/raxdns.php
